# supervisord.conf
[program:dovecot]
command = /usr/sbin/dovecot -F
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
autostart = yes
